# remotelcd.service
[Unit]
Description=Remote LCD2004 display SVXLink info
Wants=network.target
After=network.target svxlink.service
Requires=svxlink.service
PartOf=svxlink.service

StartLimitInterval=200
StartLimitBurst=100

[Service]
Type=idle
ExecStart=/usr/bin/python3 /opt/fmpoland/remotelcd/remotelcd.py
WorkingDirectory=/opt/fmpoland/remotelcd
User=root
#Nice=-10
LimitCORE=infinity
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
